# packaging/linux/codex-update-manager.service
[Unit]
Description=ChatGPT Desktop Update Manager
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/codex-update-manager daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
