# contrib/user-local-install/files/.config/systemd/user/codex-desktop-update.service
[Unit]
Description=Update ChatGPT Desktop for Linux wrapper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/codex-desktop-update --quiet
